# cerata/glaucus-s6-boot-scripts/glaucus-s6-boot-scripts/run/service/log/run
#!/usr/bin/execlineb -P
redirfd -w 2 /dev/console
redirfd -w 1 /dev/null
redirfd -rnb 0 fifo
s6-applyuidgid -u 0 -g 0
s6-log -bp T /run/log
