# cerata/glaucus-s6-boot-scripts/glaucus-s6-boot-scripts/etc/s6/run-image/service/s6-linux-init-runleveld/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
fdmove 1 3
s6-ipcserver -1c 1 -a 0700 s
s6-sudod -012
/etc/s6/scripts/runlevel
