# cerata/glaucus-s6-boot-scripts/glaucus-s6-boot-scripts/etc/s6/run-image/service/s6-svscan-log/run
#!/usr/bin/execlineb -P
redirfd -w 2 /dev/console
redirfd -w 1 /dev/console
redirfd -rnb 0 fifo
s6-setuidgid root
s6-log -bpd3 T /run/uncaught-logs
