# cerata/glaucus-s6-boot-scripts/glaucus-s6-boot-scripts/etc/s6/run/service/shutdownd/run
#!/usr/bin/execlineb -P
s6-linux-init-shutdownd -c /etc/s6
