# cerata/glaucus-s6-boot-scripts/glaucus-s6-boot-scripts/etc/s6/run-image/service/s6-linux-init-logouthookd/run
#!/usr/bin/execlineb -P
s6-ipcserver -1 -a 0700 s
s6-linux-init-logouthookd
