# cerata/glaucus-s6-boot-scripts/glaucus-s6-boot-scripts/run/service/.s6-svscan/finish
#!/usr/bin/execlineb -P
cd /
redirfd -w 1 /dev/console
fdmove -c 2 1
if { s6-svc -X /run/service/log }
if { s6-nuke -t }
if { s6-nuke -k }
if { umount -a }
if { s6-sync }
if { mount -o remount,ro / }
s6-$1 -f
