# cerata/glaucus-s6-boot-scripts/glaucus-s6-boot-scripts/run/service/.s6-svscan/finish
#!/usr/bin/execlineb -P
cd /
redirfd -w 1 /dev/console
fdmove -c 2 1
if { s6-svc -X /run/service/log }
unexport ?
s6-sync
s6-nuke -ht
s6-nuke -k
