# cerata/glaucus-s6-boot-scripts/glaucus-s6-boot-scripts/etc/s6/run-image/service/wtmpd/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
fdmove 1 3
s6-ipcserver -1 /run/utmps/.wtmpd-socket
utmps-wtmpd
